# deploy/systemd/sys-mcp-center.service
[Unit]
Description=sys-mcp-center — MCP central server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/sys-mcp-center --config /etc/sys-mcp/center.yaml
Restart=on-failure
RestartSec=5s
StandardOutput=journal
StandardError=journal
SyslogIdentifier=sys-mcp-center

NoNewPrivileges=true
ProtectSystem=strict
PrivateTmp=true

[Install]
WantedBy=multi-user.target
